# mushroom-mqtt.target
[Unit]
Description=target allowing to start/stop all instances at once
[Install]
WantedBy=multi-user.target
